# tutorsonhenry.service
[Unit]
Description=Tutors on Henry API Service
After=network.target

[Service]
# User/Group
User=nodejs
Group=nodejs

# Application Details
Type=simple
WorkingDirectory=/var/www/tutorsonhenry
ExecStart=/usr/local/bin/node /var/www/tutorsonhenry/server.js

# Environment
Environment="NODE_ENV=production"
Environment="PORT=3000"
EnvironmentFile=/var/www/tutorsonhenry/.env

# Process Management
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

# Resource Limits
LimitNOFILE=65535
LimitNPROC=512

[Install]
WantedBy=multi-user.target
